# net-misc/lemongrab/files/initd
#!/sbin/openrc-run

command="/usr/bin/lemongrab"

depend() {
    need net
}

start() {
    ebegin "Starting lemongrab"
    start-stop-daemon --start \
            --exec $command \
            --make-pidfile \
            --background \
            --user lemongrab \
            --group lemongrab \
            -1 /var/log/lemograb/stdout \
            -2 /var/log/lemograb/stderr \
            --pidfile /var/run/lemongrab.pid
    eend $?
}

stop() {
    ebegin "Stopping lemongrab"
    start-stop-daemon --stop \
            --exec $command \
            --pidfile /var/run/lemongrab.pid
    eend $?
}
