# net-misc/lemongrab/files/initd
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/bin/lemongrab"
command_user="lemongrab"
supervise_daemon_args="--stderr /var/log/lemongrab/stderr --stdout /var/log/lemongrab/stdout"
pidfile="/var/run/lemongrab.pid"

respawn_delay=10
respawn_period=300

depend() {
    use logger dns
    need net
}
